# ami/webapp.service
[Unit]
Description=CSYE6225 Web Application Service
After=network.target

[Service]
User=csye6225
Group=csye6225
WorkingDirectory=/opt/csye6225/
ExecStart=/usr/bin/node /opt/csye6225/app.js
Restart=always
RestartSec=10
EnvironmentFile=/opt/csye6225/.env
StandardOutput=append:/opt/csye6225/logs/csye6225.log
StandardError=append:/opt/csye6225/logs/csye6225.log
SyslogIdentifier=csye6225-webapp

[Install]
WantedBy=multi-user.target
